# demo-service/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:demo-service]
command=/app/demo-service
